# doc/etc-init-uwsgi.conf
description "uWSGI Python Web Server"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec /srv/grader/venv/bin/uwsgi --emperor /etc/uwsgi
