# docker/rootfs/etc/services.d/celery/run
#!/bin/execlineb -P

define user aplus
define home /srv/${user}
#define daemon /usr/bin/epmd
# Using stdout instead of stderr
fdmove -c 2 1

# (-s is undocumented magic, that seems to move _prog_ to end of _then_ or
# _else_ blocks and thus makes _define_ work)
ifthenelse -s { test -e /local/venv_${user}/bin/python3 }
{ define python3 /local/venv_${user}/bin/python3 }
{ define python3 /usr/bin/python3 }

# Use container environment
with-contenv

# user and workdir
s6-setuidgid $user
s6-env HOME=${home}
cd ${home}

# Start celery daemon
${python3} -m debugpy --listen 0.0.0.0:5679 -m celery -A aplus worker --loglevel=debug --concurrency 1 --hostname worker_main@aplus --beat --task-events --schedule /var/celery/celerybeat-schedule
